# deploy/com.secondnick.cite-cache-expansion.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.secondnick.cite-cache-expansion.plist — Task 17 recurring expansion engine.

  ┌─────────────────────────────────────────────────────────────────────────┐
  │  NOT LOADED — this is a FILE in the repo, DISARMED by design.            │
  │                                                                          │
  │  To ARM (only when Nick decides):                                        │
  │    cp deploy/com.secondnick.cite-cache-expansion.plist \                 │
  │       ~/Library/LaunchAgents/com.secondnick.cite-cache-expansion.plist   │
  │    launchctl load ~/Library/LaunchAgents/com.secondnick.cite-cache-expansion.plist
  │                                                                          │
  │  To DISARM:                                                              │
  │    launchctl unload ~/Library/LaunchAgents/com.secondnick.cite-cache-expansion.plist
  │    rm ~/Library/LaunchAgents/com.secondnick.cite-cache-expansion.plist   │
  │                                                                          │
  │  Cadence: WEEKLY, Mondays ~09:30 (Weekday=1). The starved-cell queue    │
  │  depletes faster at this cadence — runs taper to fewer/zero rows as      │
  │  thin cells fill (loop-until-dry); that's expected, not a failure.       │
  │                                                                          │
  │  MODE: `--auto` — the UNATTENDED autonomy path (ARM-B). Implies real     │
  │  headless-claude researcher (`scripts/researcher-claude.ts` → `claude    │
  │  -p`, web research), live-URL gate ON, and PUSH ON (opens a real PR via  │
  │  `gh pr create`). Because it PUSHES + opens a PR when armed, keep this    │
  │  DISARMED until Nick decides. Smoke it safely first (no PR) with:        │
  │    npx tsx scripts/run-expansion.ts --auto --dry-run --top-k=1           │
  │  Requires: a logged-in `claude` CLI + `gh` auth on this machine.         │
  └─────────────────────────────────────────────────────────────────────────┘
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.secondnick.cite-cache-expansion</string>

    <!--
      Login shell (`zsh -l`) so the user's profile PATH (npx via nvm/homebrew,
      `claude` in ~/.local/bin) is on PATH; `cd` into the repo; run the
      UNATTENDED `--auto` path. The script defaults `label` to
      `expansion-<ISO date>` when none is given.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>export PATH="$HOME/.local/bin:$PATH"; cd /Users/bignick/shared-data &amp;&amp; exec npx tsx scripts/run-expansion.ts --auto --top-k=3 --row-cap=15</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bignick/shared-data</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/bignick/shared-data/docs/pending-prs/expansion-run.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bignick/shared-data/docs/pending-prs/expansion-run.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
